# akshara-system-update.service
[Unit]
Description=Update system

[Service]
Type=simple
ExecStart=/usr/bin/akshara update
User=root

[Install]
WantedBy=multi-user.target
